# contrib/systemd/synvoid.service
[Unit]
Description=MaluWAF - Multi-Process Web Application Firewall
Documentation=https://github.com/synvoid/synvoid
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/synvoid
ExecStart=/usr/local/bin/synvoid --overseer --foreground
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5s

# Resource limits
LimitNOFILE=65535
LimitNPROC=65535

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/opt/synvoid /var/log/synvoid /var/run/synvoid

# Environment
Environment=RUST_LOG=info
Environment=SYNVOID_CONFIG_PATH=/opt/synvoid/config

# Watchdog support (optional - for systemd watchdog)
WatchdogSec=30s
NotifyAccess=all

[Install]
WantedBy=multi-user.target
